# noir-orchestrator.service
[Unit]
Description=Noir Sovereign AI Orchestrator V30.1 — 27 Pilar Evolusi
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/noir-agent/noir-vps
ExecStart=/usr/bin/python3 -u sovereign_orchestrator.py
Restart=always
RestartSec=10
StandardOutput=append:/root/noir-agent/noir-vps/orchestrator.log
StandardError=append:/root/noir-agent/noir-vps/orchestrator.log
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/noir-agent/noir-vps:/root/noir-agent/noir-ui

[Install]
WantedBy=multi-user.target
